# reference_integration/build/lola_demo/files/lola-ipc-test.service
# *******************************************************************************
# Copyright (c) 2025 Contributors to the Eclipse Foundation
#
# See the NOTICE file(s) distributed with this work for additional
# information regarding copyright ownership.
#
# This program and the accompanying materials are made available under the
# terms of the Apache License Version 2.0 which is available at
# https://www.apache.org/licenses/LICENSE-2.0
#
# SPDX-License-Identifier: Apache-2.0
# *******************************************************************************
[Unit]
Requires=bluechi-proxy@agent-qm_lola-ipc-sub.service
Requires=bluechi-proxy@agent-main_lola-ipc-pub.service
After=proxy@agent-main_lola-ipc-pub.service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/lola-ipc-test
RemainAfterExit=yes
